# broadlinkac-to-mqtt/rootfs/etc/services.d/ac2mqtt/run
#!/usr/bin/with-contenv bashio


bashio::log.info "Starting ac2mqtt..."
exec \
	python /app/monitor.py --config /data/config.yml --data_dir /data
